# rootfs/etc/services.d/hostmaster/run
#!/usr/bin/with-contenv bash
bash ./hostmaster.sh
exec gosu aegir drush @hostmaster hosting-queued
